# config/usr/lib/systemd/user/vdr-sxfe.service
[Unit]
Description=start vdr-sxfe as vdr frontend

[Service]
Environment=OPTIONS="-E localhost"
EnvironmentFile=-%h/.config/yavdr-frontend/xineliboutput.conf
Type=simple
ExecStart=/usr/bin/vdr-sxfe $OPTIONS
Restart=on-failure
